# package/zustomers/dg301/ELION/fs/etc/hotplug.d/iface/40-nameserver

[ ifup = "$ACTION" ] && {
	case "$INTERFACE" in
		internet|mobile) nvram set dnsmasq_populate=1 ;;
	esac
}

