# package/zustomers/dg301/ELION/fs/etc/hotplug.d/iface/45-neighlogger

[ ifup = "$ACTION" ] {
        local islan="$(uci -q get network.$INTERFACE.is_lan)"
        if [ "$INTERFACE" != "loopback" ] && [ "$islan" == "1" ]; then
                pidof neighlogger >/dev/null 2>&1 || /sbin/neighlogger &
        fi
}
