# package/iopsys-base-files/iopsys/base-files-iopsys/etc/hotplug.d/iface/60-wanuptest

[ ifup = "$ACTION" -o ifdown = "$ACTION" ] && {
	local islan="$(uci -q get network.$INTERFACE.is_lan)"
	local proto="$(uci -q get network.$INTERFACE.proto)"
	local ifname="$(uci -q get network.$INTERFACE.ifname)"
	if [ "$islan" != "1" ] && [ "$proto" != "none" ] && [ "${ifname:0:1}" != "@" ]; then
		killall testnet >/dev/null 2>&1
		/sbin/testnet $INTERFACE &
	fi
}
