# target/linux/brcm63xx_ref_nand/base-files/etc/hotplug.d/iface/60-wanuptest

[ ifup = "$ACTION" -o ifdown = "$ACTION" ] && {
	local islan="$(uci -q get network.$INTERFACE.is_lan)"
	local proto="$(uci -q get network.$INTERFACE.proto)"
	if [ "$islan" != "1" ] && [ "$proto" != "none" ]; then
		killall testnet >/dev/null 2>&1
		/sbin/testnet &
	fi
}
