# target/linux/brcm63xx_ref_nand/base-files/etc/hotplug.d/iface/60-wanuptest

[ ifup = "$ACTION" -o ifdown = "$ACTION" ] && {
	if [ "$(uci -q get network.$INTERFACE.is_lan)" != "1" ]; then
		pidof testnet >/dev/null || /sbin/testnet &
	fi
}
