# target/linux/brcm63xx_ref_nand/base-files/etc/hotplug.d/iface/70-shiftrange

[ ifup = "$ACTION" ] && {
	if [ "$(uci -q get network.$INTERFACE.is_lan)" != "1" ]; then
		. /lib/functions/network.sh
		network_get_ipaddr wan_ip wan

		shift_range $wan_ip && {
			/etc/init.d/network reload
			/etc/init.d/firewall reload
			/etc/init.d/dnsmasq restart
			/etc/init.d/mcpd restart
			/etc/init.d/dropbear restart
			/etc/init.d/samba restart
			/etc/init.d/domain_name_access restart
		}
	fi
}

