# target/linux/brcm63xx_ref_nand/base-files/etc/hotplug.d/switch/10-switch
. /lib/network/config.sh

addif() {
	# check if wan had a event and start all related services
	wan=`uci -q get layer2_interface_ethernet.Wan.baseifname`
	if [ "$INTERFACE" == "$wan" ]; then
		/etc/init.d/layer2_interface_ethernet start
		/etc/init.d/layer2_interface_vlan start
	fi
}

delif() {
	##remove ethernet wan and nuke any vlans associated
	wan=`uci -q get layer2_interface_ethernet.Wan.baseifname`
	if [ "$INTERFACE" == "$wan"  ]; then
		/etc/init.d/layer2_interface_ethernet stop
	fi
}

gigatest() {
	local speed=$(ethctl $INTERFACE media-type 2>&1 | awk '{if (NR == 2) print $6}')
	case "$speed" in
		1000*) return 1 ;;
		*) return 0 ;;
	esac
}

ethtest () {
	local landev=$(db get hw.board.ethernetLanPorts)
	local ledontest
	for dev in $landev; do
		ledontest=$(cat /sys/class/net/$dev/operstate)
		if [ "$ledontest" == "up" ]; then
			return 1
		fi
	done
	return 0
}

dsltest() {
	if cat /var/state/layer2_interface 2>/dev/null | grep 'adsl\|vdsl' | grep up; then
		return 1
	else
		return 0
	fi
}

case "$ACTION" in
	add|register)
		case "$PHYSDEVDRIVER" in
			natsemi) sleep 1 ;;
		esac
		addif
		local interfname=$(interfacename $INTERFACE)
		local operstate=$(cat /sys/class/net/$INTERFACE/operstate)
		if [ "$operstate" == "up" ]; then
			case "$interfname" in
				GbE*)
					ledctl gbe off
					$(gigatest) || ledctl gbe alert 
					ledctl gbe notice
					ledctl lan on
				;;
				LAN*)
					ledctl lan on
				;;
				WAN*)
					ledctl wan off
					$(gigatest) || ledctl wan alert
					ledctl wan notice
					ledctl broadband on
				;;
			esac
		fi		
	;;
	remove|unregister)
		local interfname=$(interfacename $INTERFACE)
		local operstate=$(cat /sys/class/net/$INTERFACE/operstate)
		if [ "$operstate" == "down" ]; then
			case "$interfname" in
				GbE*)
					ledctl gbe off
					$(ethtest) && ledctl lan off
				;;
				LAN*)
					$(ethtest) && ledctl lan off
				;;
				WAN*)
					ledctl wan off
					$(dsltest) && ledctl broadband off
				;;
			esac
		fi
		delif
	;;
esac

