# target/linux/brcm63xx_ref_nand/base-files/etc/hotplug.d/switch/10-switch

include /lib/network

addif() {
	# check if wan had a event and start all related services
	wan=`uci -q get layer2_interface_ethernet.Wan.baseifname`
	if [ "$INTERFACE" == "$wan" ]; then
		/etc/init.d/layer2_interface_ethernet start
		/etc/init.d/layer2_interface_vlan start
	fi
}

delif() {
	##remove ethernet wan and nuke any vlans associated
	wan=`uci -q get layer2_interface_ethernet.Wan.baseifname`
	if [ "$INTERFACE" == "$wan"  ]; then
		/etc/init.d/layer2_interface_ethernet stop
	fi
}

ethtest () {
	local landev=$(db get hw.board.ethernetPortOrder)
	local ledontest
	for dev in $landev; do
		ledontest=$(cat /sys/class/net/"$dev"/operstate)
		if [ "$ledontest" == "up" ]; then
			return 1
		fi
	done
	return 0
}

case "$ACTION" in
	add|register)
		case "$PHYSDEVDRIVER" in
			natsemi) sleep 1 ;;
		esac
		local interfname=$(interfacename $INTERFACE)
		local speed=$(interfacespeed $INTERFACE)
		ledon=$(cat /sys/class/net/"$INTERFACE"/carrier)$interfname$speed
        operstate=$(cat /sys/class/net/"$INTERFACE"/operstate)
        if [ "$operstate" == "up" ]; then
            ubus call led.lan set '{"state":"ok"}' &
        fi

		addif
	;;
	remove|unregister)
		local interfname=$(interfacename $INTERFACE)
		ledon=$(cat /sys/class/net/"$INTERFACE"/carrier)$interfname
        operstate=$(cat /sys/class/net/"$INTERFACE"/operstate)
        if [ "$operstate" == "down" ]; then
            if $(ethtest); then
                ubus call led.lan set '{"state":"off"}' &
            fi
        fi


		delif
	;;
esac
