# package/customers/dg301/ELION/fs/etc/hotplug.d/iface/55-asterisk_log
[ ifup = "$ACTION" ] && {
        if [ "$INTERFACE" == "internet" ]; then
                /bin/echo "syslog.local0 => notice,warning,error" >> /etc/asterisk_templates/logger.conf.TEMPLATE
                /bin/echo "syslog.local0 => notice,warning,error" >> /etc/asterisk/logger.conf
                asterisk -rx "logger reload"
       fi
}

[ ifdown = "$ACTION" ] && {
       if [ "$INTERFACE" == "internet" ]; then
                /bin/sed -i "/^syslog/d" /etc/asterisk_templates/logger.conf.TEMPLATE
                /bin/sed -i "/^syslog/d" /etc/asterisk/logger.conf
                /etc/init.d/asterisk restart &
      fi
}
