# package/inteno/questd/files/etc/init.d/quest
#!/bin/sh /etc/rc.common

START=14

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

start() {
        service_start /sbin/questd
}

stop() {
        service_stop /sbin/questd
}

