# package/zeta/basefiles/iopsys/base-files-iopsys/etc/hotplug.d/dsl/10-dsl
include /lib/network

addif() {
	#find all vlan configurations for this interface and set them up as well
	/etc/init.d/layer2_interface_vlan start
}

start() {
	case "$ACTION" in
		add|register)
			addif
		;;
		remove|unregister)                                               
		;;
	esac
}

if [ "ADSL" == "$INTERFACE" -o "VDSL2" == "$INTERFACE" ]; then
	if [ "ADSL" == "$INTERFACE" ]; then 
		if [ "$ACTION" == "remove" ]; then
			uci_toggle_state layer2_interface adsl device down
			/etc/init.d/layer2_interface_adsl stop
		else
			uci_toggle_state layer2_interface adsl device up
			/etc/init.d/layer2_interface_adsl start
		fi
	fi
	if [ "VDSL2" == "$INTERFACE" ]; then
		if [ "$ACTION" == "remove" ]; then
			uci_toggle_state layer2_interface vdsl device down
			/etc/init.d/layer2_interface_vdsl  stop
		else
			uci_toggle_state layer2_interface vdsl device up
			/etc/init.d/layer2_interface_vdsl start
		fi
	fi
fi  

start
