# package/base-files/iopsys/etc/hotplug.d/switch/20-arpscan
. /lib/network/config.sh

case "$ACTION" in
	add|register)
		netcheck -i $(get_network_of $INTERFACE) -q 10
	;;
	remove|unregister)
	;;
esac

