# target/linux/brcm63xx_ref_nand/base-files/etc/init.d/enviroment
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=11

handle_threads()
{
	local pid threads

	threads="wl0-kthrd wl0-tmr-kthrd"
	if wlctl -i wl1 ap >/dev/null 2>&1; then
		threads="$threads wl1-kthrd wl1-tmr-kthrd"
	fi

	# initialize kernel schedule priorities
	smd -s

	# set priority for given threads
        for thread in $threads; do
                pid=$(ps | grep $thread | head -1 | awk '{print$1}')
		[ -n "$pid" ] && chrt -rp 5 $pid
        done
}

handle_configs()
{
	local BMAC=$(cat /proc/nvram/BaseMacAddr | tr '[a-z]' '[A-Z]')
	BMAC=${BMAC// /}
	local MAC=$(printf "%X\n" $((0x$BMAC)))
	local BSSID=$(printf "%X\n" $((0x$BMAC + 2)))
	local WPAKEY=$(cat /proc/nvram/WpaKey)
	local SERIALNR=$(cat /proc/nvram/SerialNumber)
	local HWVER=$(db get hw.board.hardwareVersion)
	local RMODEL=$(db get hw.board.routerModel)
	local DESKEY=$(db get hw.board.desKey)

	local oid=${BMAC:0:6}
	local mac=$BMAC
	local mac2=${MAC:8:2}
	local mac4=${MAC:6:4}
	local mac6=${MAC:4:6}
	local bssid=$BSSID
	local bssid2=${BSSID:8:2}
	local bssid4=${BSSID:6:4}
	local bssid6=${BSSID:4:6}
	local wpakey="${WPAKEY:-1234567890}"
	local hardwareid=$HWVER-$(echo $RMODEL | sed -r 's;.+-(.+);\1;')

	local configs="passwords network wireless system cwmp provisioning"

	for config in $configs; do
		if [ -f /etc/config/$config ]; then
			sed -i "s/\$MAC6/$mac6/g" /etc/config/$config
			sed -i "s/\$MAC4/$mac4/g" /etc/config/$config
			sed -i "s/\$MAC2/$mac2/g" /etc/config/$config
			sed -i "s/\$MAC/$mac/g" /etc/config/$config
			sed -i "s/\$BSSID6/$bssid6/g" /etc/config/$config
			sed -i "s/\$BSSID4/$bssid4/g" /etc/config/$config
			sed -i "s/\$BSSID2/$bssid2/g" /etc/config/$config
			sed -i "s/\$BSSID/$bssid/g" /etc/config/$config
			sed -i "s/\$WPAKEY/$wpakey/g" /etc/config/$config
			sed -i "s/\$DESKEY/$DESKEY/g" /etc/config/$config
			sed -i "s/\$SER/$SERIALNR/g" /etc/config/$config
			sed -i "s/\$OUI/$oid/g" /etc/config/$config
			sed -i "s/\$HARDWAREID/$hardwareid/g" /etc/config/$config
			[ "$config" == "network" ] && sed -i "s/multiwan/anywan/g" /etc/config/$config
		fi
	done
}

start()
{
	handle_threads
	handle_configs
}

restart()
{
	handle_configs
}

