# target/linux/brcm63xx_ref_nand/base-files/etc/init.d/makedevs
#!/bin/sh /etc/rc.common

START=01

start() {
	echo "Copying device files from /lib/dev to /dev ...again"
	cp -a /lib/dev/* /dev
}
