# target/linux/brcm63xx_ref_nand/base-files/etc/init.d/passwords
#!/bin/sh /etc/rc.common

. /lib/functions.sh

START=12

cwmp_passwordcache() {
        [ -f /usr/sbin/freecwmp ] && {
                ps -w | grep freecwmp | grep InternetGatewayDevice.X_INTENO_SE_LoginCfg >/dev/null || freecwmp get cache InternetGatewayDevice.X_INTENO_SE_LoginCfg. &
        }
}

migrate_shadow() {
	/bin/sh /rom/etc/uci-defaults/migrate-shadow
}

set_password() {
        local usertype="$1"
        local password
	config_get password $usertype password

	[ -n "$password" ] && (echo $password; sleep 1; echo $password) | passwd $usertype >/dev/null 2>&1

	uci delete passwords.$usertype.password >/dev/null 2>&1
}

start() {
	migrate_shadow
	config_load passwords
	config_foreach set_password usertype
	uci commit passwords
	cwmp_passwordcache
}
