# target/linux/brcm63xx_ref_nand/base-files/etc/init.d/passwords
#!/bin/sh /etc/rc.common

. /lib/functions.sh

START=12

set_password() {
        local usertype="$1"
        local password
	config_get password $usertype password

	[ -n "$password" ] && (echo $password; sleep 1; echo $password) | passwd $usertype >/dev/null 2>&1

	uci delete passwords.$usertype.password >/dev/null 2>&1
}

start() {
	config_load passwords
	config_foreach set_password usertype
	uci commit passwords
}
