# target/linux/brcm63xx_ref_nand/base-files/etc/init.d/swmdk
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2009 OpenWrt.org
#

START=41

start() {
        rm /lib/libc.so.0
        ln -s /usr/lib/libc.so.0 /lib/libc.so.0
        /usr/sbin/swmdk &
        sleep 3
        rm /lib/libc.so.0
        ln -s /lib/libuClibc-0.9.33.2.so /lib/libc.so.0
}
