# package/zustomers/dg301/NGT/fs/etc/init.d/ngtwan
#!/bin/sh /etc/rc.common

. /lib/functions.sh

START=99

start() {
	/sbin/ngtwan &
}
