# package/ledmngr/files/etc/init.d/ledmngr
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2009 OpenWrt.org
#

START=39

start() {
        /sbin/ledmngr
}

