# package/customers/cg300/TELIASAFETY/fs/etc/hotplug.d/iface/90-m2msupervision

[ ifup = "$ACTION" -a "$INTERFACE" = "mobile" ] && {
	ps | grep "m2msupervision" | grep -v "grep" >/dev/null || m2msupervision &
}

