# target/linux/brcm63xx_ref_nand/base-files/files/etc/init.d/swmdk
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2009 OpenWrt.org
#

START=41

start() {
	export LD_LIBRARY_PATH=/usr/lib
        /usr/sbin/swmdk &
}
