# target/linux/brcm63xx_ref_nand/base-files/etc/hotplug.d/button/10-button
#!/bin/sh

. /usr/share/libubox/jshn.sh

case "$ACTION" in
        add|register)
        [ "resetbutton" == "$INTERFACE" ] && {
                ubus call leds set '{ "state": "test"}'
                printf "\n\n*** Restore to default settings ***\n\n" > /dev/console
		mtd -r erase rootfs_data 2>/dev/null
        }
        [ "ecobutton" == "$INTERFACE" ] && {
		if ! pidof ecobutton; then
			/sbin/ecobutton &
		fi
		local buttoncomb=""
		[ -f "/var/buttoncomb" ] && buttoncomb="$(cat /var/buttoncomb)"
		echo  $buttoncomb"eco" > /var/buttoncomb
        }
        [ "wpsbutton" == "$INTERFACE" ] && {
		        printf "\n\nWPS button is pressed\n\n" > /dev/console
        }
        [ "dectbutton" == "$INTERFACE" ] && {
                printf "\n\nRegistering DECT\n\n" > /dev/console
                local parsed=$(ubus call led.dect status)
                json_load "$parsed"
                json_get_var ledstate state

                [ "$ledstate" == "off" ] && {
                        ubus call led.dect set '{ "state": "notice"}'
                } || {
                        ubus call led.dect set '{ "state": "off"}'
                }
        }
        [ "infobutton" == "$INTERFACE" ] && {
                printf "\n\nINFO button is pressed\n\n" > /dev/console
                local parsed=$(ubus call leds status)
                json_load "$parsed"
                json_get_var ledstate state

                [ "$ledstate" == "normal" ] && {
                        ubus call leds set '{ "state": "info"}'
                } || {
                        ubus call leds set '{ "state": "normal"}'
                }
        }
        ;;
        remove|unregister)
        ;;
esac
