# target/linux/brcm63xx_ref_nand/base-files/etc/hotplug.d/button/10-button
#!/bin/sh

. /lib/functions.sh
. /usr/share/libubox/jshn.sh

wpsbutton_action () {
        killall -USR2 wps_monitor
}

dectbutton_action () {
        local parsed=$(ubus call led.dect status)
        json_load "$parsed"
        json_get_var ledstate state
	/usr/bin/dect -r > /dev/null

        [ "$ledstate" == "off" ] && {
                ubus call led.dect set '{ "state": "notice"}'
        } || {
                ubus call led.dect set '{ "state": "off"}'
        }
}

case "$ACTION" in
        add|register)
        [ "resetbutton" == "$INTERFACE" ] && {
                ubus call leds set '{ "state": "test"}'
                printf "\n\n*** Restore to default settings ***\n\n" > /dev/console
		rm -rf /rom/overlay/*
		rm -rf /rom/etc/adsl/adsl_phy.bin
		reboot
        }
        [ "ecobutton" == "$INTERFACE" ] && {
		if ! pidof ecobutton; then
			/sbin/ecobutton &
		fi
		local buttoncomb=""
		[ -f "/var/buttoncomb" ] && buttoncomb="$(cat /var/buttoncomb)"
		echo  $buttoncomb"eco" > /var/buttoncomb
        }
        [ "wpsbutton" == "$INTERFACE" ] && {
	        printf "\n\nWPS button is pressed\n\n" > /dev/console
		wpsbutton_action
        }
        [ "dectbutton" == "$INTERFACE" ] && {
                printf "\n\nRegistering DECT\n\n" > /dev/console
		dectbutton_action
        }
        [ "infobutton" == "$INTERFACE" ] && {
                printf "\n\nINFO button is pressed\n\n" > /dev/console
                local parsed=$(ubus call leds status)
                json_load "$parsed"
                json_get_var ledstate state

                [ "$ledstate" == "normal" ] && {
                        ubus call leds set '{ "state": "info"}'
                        ubus call led.eco set '{ "state": "notice"}'
                } || {
                        ubus call leds set '{ "state": "normal"}'
                        ubus call led.eco set '{ "state": "off"}'
                }
        }
	[ "communicationbutton" == "$INTERFACE" ] && {
		printf "\n\nCommunication button pressed - not implemented!\n\n" > /dev/console
	}
	[ "servicebutton" == "$INTERFACE" ] && {
		printf "\n\nService button pressed - not implemented!\n\n" > /dev/console
	}
	[ "pairbutton" == "$INTERFACE" ] && {
		printf "\n\nPair button pressed\n\n" > /dev/console
		wpsbutton_action
		dectbutton_action
	}
	[ "cancelbutton" == "$INTERFACE" ] && {
		printf "\n\nCancel button pressed - not implemented!\n\n" > /dev/console
	}
	[ "touch_near" == "$INTERFACE" ] && {
	        ubus call leds proximity '{ "timeout" : "15", "light-all" : "2" }'
	}
	[ "touch_far" == "$INTERFACE" ] && {
	        ubus call leds proximity '{ "timeout" : "300" }'
	}
        ;;
        remove|unregister)
        ;;
esac
