# target/linux/brcm63xx_ref_nand/base-files/etc/hotplug.d/button/10-button
#!/bin/sh

. /lib/functions.sh
. /usr/share/libubox/jshn.sh

case "$ACTION" in
        add|register)
        [ "resetbutton" == "$INTERFACE" ] && {
                ubus call leds set '{ "state": "test"}'
                printf "\n\n*** Restore to default settings ***\n\n" > /dev/console
                umount "$(find_mtd_part rootfs_data)"
                echo -e "\xde\xad\xc0\xde" | mtd -x -qq write - rootfs_data
		reboot
        }
        [ "ecobutton" == "$INTERFACE" ] && {
		if ! pidof ecobutton; then
			/sbin/ecobutton &
		fi
		local buttoncomb=""
		[ -f "/var/buttoncomb" ] && buttoncomb="$(cat /var/buttoncomb)"
		echo  $buttoncomb"eco" > /var/buttoncomb
        }
        [ "wpsbutton" == "$INTERFACE" ] && {
		        printf "\n\nWPS button is pressed\n\n" > /dev/console
		        killall -USR2 wps_monitor
        }
        [ "dectbutton" == "$INTERFACE" ] && {
                printf "\n\nRegistering DECT\n\n" > /dev/console
                local parsed=$(ubus call led.dect status)
                json_load "$parsed"
                json_get_var ledstate state
		/usr/bin/dect -r > /dev/null

                [ "$ledstate" == "off" ] && {
                        ubus call led.dect set '{ "state": "notice"}'
                } || {
                        ubus call led.dect set '{ "state": "off"}'
                }
        }
        [ "infobutton" == "$INTERFACE" ] && {
                printf "\n\nINFO button is pressed\n\n" > /dev/console
                local parsed=$(ubus call leds status)
                json_load "$parsed"
                json_get_var ledstate state

                [ "$ledstate" == "normal" ] && {
                        ubus call leds set '{ "state": "info"}'
                } || {
                        ubus call leds set '{ "state": "normal"}'
                }
        }
	[ "touch_near" == "$INTERFACE" ] && {
	        ubus call leds proximity '{ "timeout" : "15", "light-all" : "2" }'
	}
	[ "touch_far" == "$INTERFACE" ] && {
	        ubus call leds proximity '{ "timeout" : "10" }'
	}
        ;;
        remove|unregister)
        ;;
esac
