# package/zustomers/vox25/VODAFONE/fs/etc/init.d/vodafone
#!/bin/sh /etc/rc.common

START=99

boot() {
        /etc/init.d/samba reload
}
