# package/zustomers/dg301/ELION/fs/etc/crontabs/root
*/5 * * * * /usr/sbin/logrotate /etc/logrotate.d/messages -s /var/log/rotate.status
0 0 * * * /usr/sbin/logrotate -f /etc/logrotate.d/messages -s /var/log/rotate.status
