# package/netifd/files/etc/hotplug.d/iface/60-wanuptest

[ ifup = "$ACTION" -o ifdown = "$ACTION" ] && {
	/sbin/testnet &
}
