# quarantine/old-rebase-on-master/package/peripheral_manager/files/etc/init.d/peripheral_manager
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2009 OpenWrt.org
#

START=19

start() {
        /sbin/peripheral_manager
        sleep 1
        /bin/ubus call led.status set '{"state":"ok"}' &
}

stop() {
        killall -9 peripheral_manager
}
