# package/customers/vg50/TTT2/fs/etc/crontabs/root
0 03 * * *     /sbin/iup > /dev/null 2>&1
