# target/linux/brcm63xx_ref_nand/base-files/etc/hotplug.d/usb/10-usb
#!/bin/sh

# Copyright (C) 2009 OpenWrt.org

case "$ACTION" in
	add)
		if [ -d /sys/bus/usb/devices/1-1 -o -d /sys/bus/usb/devices/1-2 -o -d /sys/bus/usb/devices/2-1 -o -d /sys/bus/usb/devices/2-2 ] ; then
		  /usr/sbin/ledctl USB on
		fi
		;;
	remove)
		if [ ! -d /sys/bus/usb/devices/1-1 -a ! -d /sys/bus/usb/devices/1-2 -a ! -d /sys/bus/usb/devices/2-1 -a ! -d /sys/bus/usb/devices/2-2 ]; then
		  /usr/sbin/ledctl USB off
		fi
		;;
esac	
		
